# matrix-switch.service
[Unit]
Description=OREI Matrix Switch Web Interface
After=network-online.target
Wants=network-online.target

[Service]
User=matrix
Group=matrix
WorkingDirectory=/opt/matrix/orei-matrix-controller
ExecStart=/opt/matrix/.venv/bin/gunicorn app:app \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --worker-class gthread \
    --threads 4 \
    --timeout 120 \
    --access-logfile /var/log/matrix-switch/access.log \
    --error-logfile /var/log/matrix-switch/error.log
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
